# miniconda3/pkgs/conda-25.11.1-py313hca03da5_0/etc/profile.d/conda.csh
setenv CONDA_EXE "/opt/miniconda3/conda-bld/conda_1765567809749/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_pl/bin/conda";
setenv _CONDA_EXE "/opt/miniconda3/conda-bld/conda_1765567809749/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_pl/bin/conda";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "/opt/miniconda3/conda-bld/conda_1765567809749/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_pl/bin/python";
setenv _CONDA_ROOT "/opt/miniconda3/conda-bld/conda_1765567809749/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_pl";

source "/opt/miniconda3/conda-bld/conda_1765567809749/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_pl/lib/python3.13/site-packages/conda/shell/etc/profile.d/conda.csh"